# llama-updater.service
[Unit]
Description=Check for and apply llama.cpp updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-llama-server
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/llama-server /var/lib/llama-updater /var/cache/llama-updater
DevicePolicy=auto

# Grant access to GPU devices
DeviceAllow=/dev/dri rw
BindReadOnlyPaths=/dev/dri
